# scripts/idleshutdown.service
[Unit]
Description=IdleShutdown Agent - VM Idle Monitoring Service
Documentation=https://github.com/sricharan-11/vm-idle-shutdown
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/idleshutdown --config /etc/idleshutdown/config.ini
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=IdleShutdown

# Security hardening
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/etc/idleshutdown
NoNewPrivileges=no

[Install]
WantedBy=multi-user.target
